# project-b/contrib/systemd/npu-status.timer
[Unit]
Description=Weekly XDNA 1 status probe (drives npu-status.service)

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h
Unit=npu-status.service

[Install]
WantedBy=timers.target
